# hier/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:node]
command=/usr/bin/node /opt/app/build/server

[program:nginx]
command=/etc/init.d/nginx start
